# packaging/satur8-daemon.service
[Unit]
Description=Satur8 focus daemon (per-game digital vibrance)
Documentation=https://github.com/NtrpyDev/satur8
# Only meaningful inside a graphical session.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/satur8-daemon
Restart=on-failure
RestartSec=2
# Event-driven: it sleeps until the KWin script forwards a focus change.

[Install]
WantedBy=graphical-session.target
